# supervisord.conf
[supervisord] 
nodaemon=true

[benutzer:fhem]
command=chown -R /opt/fhem

[program:fhem] 
command=/etc/init.d/fhem start
